# q-system/.q-system/scripts/com.kipi.openloops.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Durable open-loops reminder. Install: cp to ~/Library/LaunchAgents/ then
     launchctl load it. Fires Mon + Thu 9:23am; pops a desktop notification when
     something is tagged [needs you]. Disable: launchctl unload the plist. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kipi.openloops</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/assafkipnis/projects/kipi-system/q-system/.q-system/scripts/open-loops-notify.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>23</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>23</integer></dict>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardErrorPath</key>
  <string>/tmp/kipi-openloops.err</string>
</dict>
</plist>
